# systemd/odoo-mcp.service
[Unit]
Description=Odoo MCP Server — Heinzbot Integration
After=network.target odoo.service

[Service]
Type=simple
User=odoo
WorkingDirectory=/opt/odoo-mcp
EnvironmentFile=/opt/odoo-mcp/.env
ExecStart=/opt/odoo-mcp/.venv/bin/uvicorn src.main:app --host 0.0.0.0 --port 8001
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
